# dist/arch/arsse.service
[Unit]
Description=The Arsse newsfeed management service
Documentation=https://thearsse.com/manual/
Requires=arsse-fetch.service
BindsTo=php-fpm.service
After=php-fpm.service

[Service]
Type=oneshot
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
